# deploy/systemd/pedestrian-single-loop.service.example
[Unit]
Description=Pedestrian Line Counter (single-loop live)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pedline
Group=pedline
WorkingDirectory=/opt/pedestrian-line
EnvironmentFile=/etc/pedline/single-loop.env
ExecStart=/opt/pedestrian-line/scripts/run_single_loop_live.sh
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pedline-single-loop
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
