# deploy/systemd/pedestrian-edge-service.service.example
[Unit]
Description=Pedestrian Line Counter (edge FastAPI service)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pedline
Group=pedline
WorkingDirectory=/opt/pedestrian_line_counter
EnvironmentFile=/etc/vehicle_count/edge_service.env
ExecStart=/opt/pedestrian_line_counter/scripts/run_edge_service.sh
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pedline-edge-service
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
